# config_files/init/uwsgi.conf
description "uwsgi tiny instance"
start on runlevel [2345]
stop on runlevel [06]

respawn

exec uwsgi --ini /home/master/wifi_avis/wifi_avis_uwsgi.ini
